# 01-wazuh/local_rules/103-scheduled-task-outofhours.xml
<!-- 103-scheduled-task-outofhours  -->

<group name="local,scheduled_task">
  <rule id="100103" level="8">
    <field name="event_id">4698</field>
    <description>Scheduled task created (EventID 4698) - flagged if created outside business hours</description>
    <options>no_full_log</options>
    <group>persistence,tampering</group>
    <mitre>
      <id>T1053</id>
    </mitre>
  </rule>
</group>
